# reinstall/files/rc.local
#!/usr/bin/env bash

set -euo pipefail

disk=$(ls /sys/class/block | sort -n | head -n 1)
growpart "/dev/${disk}" 2
resize2fs "$(findmnt / -o SOURCE -n)"
